# src/main/centos/etc/rc.local
#!/bin/bash
set -xev
source /etc/profile

port=$1
if [ ! -n "$port" ]
then
    port=80      
fi

echo "use $port as a default" 

echo "解压war包......"
java -jar /usr/local/ngrinder.war --port $port






